# ops/systemd/footbag.service
# =============================================================================
# /etc/systemd/system/footbag.service
#
# Manages the Footbag Platform Docker Compose stack on the Lightsail host.
# Install: sudo cp ops/systemd/footbag.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable footbag
#          sudo systemctl start footbag
#
# Runs the compose stack defined by docker/docker-compose.yml + prod overrides.
# Host bootstrap (Docker CE install, /srv/footbag layout, /srv/footbag/env
# creation) must be completed before enabling this unit.
# =============================================================================

[Unit]
Description=Footbag Platform (Docker Compose)
Documentation=https://github.com/davidleberknight/footbag-platform
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Working directory is the project root on the host
WorkingDirectory=/srv/footbag

# Environment file provides PUBLIC_BASE_URL, FOOTBAG_DB_PATH, etc.
# Created manually on the host during bootstrap (root-owned, 0600).
EnvironmentFile=/srv/footbag/env

# Images must already be present on the host before this unit starts (built
# on the operator workstation by deploy-code.sh / deploy-rebuild.sh, then
# transferred via docker save | docker load). The unit refuses to build at
# start time (--no-build): a build at boot can OOM small instances and wedge
# the host before sshd is responsive. If an image is missing, the unit fails
# fast with a clear error rather than silently triggering a build.

# Start the stack; --remove-orphans cleans up containers from old service defs.
# --no-build forbids boot-time image builds (see comment above).
ExecStart=/usr/bin/docker compose \
    --env-file /srv/footbag/env \
    -f docker/docker-compose.yml \
    -f docker/docker-compose.prod.yml \
    up --detach --remove-orphans --no-build

# Graceful stop — waits up to 30s for in-flight requests to complete
ExecStop=/usr/bin/docker compose \
    --env-file /srv/footbag/env \
    -f docker/docker-compose.yml \
    -f docker/docker-compose.prod.yml \
    down --timeout 30

# Reload without downtime (nginx reload + rolling restart for web).
# --no-build for the same reason as ExecStart.
ExecReload=/usr/bin/docker compose \
    --env-file /srv/footbag/env \
    -f docker/docker-compose.yml \
    -f docker/docker-compose.prod.yml \
    up --detach --remove-orphans --no-build

Restart=on-failure
RestartSec=10s
TimeoutStartSec=120s
TimeoutStopSec=60s

[Install]
WantedBy=multi-user.target
